# ops/launchd/com.folloze.codex-skills-sync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.folloze.codex-skills-sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd __REPO_ROOT__ &amp;&amp; git pull --ff-only &amp;&amp; python3 scripts/sync_codex_skills.py --overwrite</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.folloze.codex-skills-sync.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.folloze.codex-skills-sync.err</string>
  </dict>
</plist>
